# fpp-ha-stats/services/server/run
#!/usr/bin/with-contenv bash
cd /app/server
exec node server.js
