# fpp-ha-stats/services/collector/run
#!/usr/bin/with-contenv bash
cd /app/statsCollector
exec node collector.js
