# fpp-ha-stats/services/fpp-stats/run
#!/bin/sh
# Launch processes manually
cd /app/server
node server.js &
cd /app/statsCollector
node collector.js &
wait